# supervisor/conf.d/sshd.conf
[program:sshd]
redirect_stderr=true
priority=500
command=/usr/sbin/sshd -D
